# services/mpd/run
#! /usr/bin/env sh

cd
exec mpd --no-daemon
