# services/transmission/run
#! /usr/bin/env sh

cd
exec transmission-daemon -f
