# services/mailwatch/run
#! /bin/sh

cd
exec mailwatch ~/Maildir
