# services/disk_indicator/run
#! /usr/bin/env sh

cd
exec disk_indicator -x