# services/xmobar/run
#! /bin/sh

# exec xmobar ~/cfg/my/xmobarrc
