# services/bitlbee/run
#! /usr/bin/env sh

exec bitlbee -c ../../bitlbee/bitlbee.conf -n
