# services/compton/run
#! /usr/bin/env sh

# cd
# exec compton
