# services/xscreensaver/run
#! /usr/bin/env sh

cd
exec xscreensaver -no-splash
